# oliveheart.service
[Unit]
Description=oliveheart redmine configuration service

[Service]
ExecStart=/opt/oliveheart/oliveheart.sh

[Install]
WantedBy=multi-user.target